# dock-helper.path
[Unit]
Description=Dock helper service (path activation)

[Path]
PathExists=/dev/thinkpad/dock

[Install]
WantedBy=paths.target
