# panels.target
[Unit]
Description=Panels
Documentation=man:systemd.special(7)
StopWhenUnneeded=yes
