# panels.target
[Unit]
Description=Panels
Documentation=man:systemd.special(7)
Requires=hlwm-panels.service
StopWhenUnneeded=yes
