# unlock.target
[Unit]
Description=Services which should run when session is unlocked
Conflicts=lock.target
StopWhenUnneeded=yes
